# data/ginn.conf
description "GINN multi-touch gesture reinjector"

start on started dbus and (xsession SESSION=ubuntu)
stop on desktop-end

exec /usr/bin/ginn
